# data_load_service/supervisord.conf
[supervisord]
nodaemon=true

[program:sheduller]
command=python /app/db_module.py
autorestart=true

[program:uvicorn]
command=uvicorn run:app --host 0.0.0.0 --port 8001
autorestart=true
